# nasmail/runit/dovecot/run
#!/bin/sh
[ -f /run/dovecot/master.pid ] && rm /run/dovecot/master.pid
exec /usr/sbin/dovecot -F
